# files/system/etc/udev/rules.d/98-thunderbolt-rescan.rules
ACTION=="change", SUBSYSTEM=="thunderbolt", ATTR{authorized}=="1", RUN+="/usr/local/bin/thunderbolt-rescan.sh"